# lib/supervisor/conf.d/minecraft.conf
[program:minecraft]
startsecs=5
directory=/mnt/minecraft/
command=java -Xmx1024M -Xms1024M -jar craftbukkit.jar -o true
stderr_logfile = /mnt/minecraft/logs/craftbukkit-stderr.log
stdout_logfile = /mnt/minecraft/logs/craftbukkit-stdout.log
